# example/mnt-gvol0.mount
[Unit]

[Mount]
What=127.0.0.1:/gvol0
Where=/mnt/gvol0
Type=glusterfs-docker

[Install]
WantedBy=multi-user.target
