# example/mnt-gvol0.automount
[Unit]
[Automount]
Where=/mnt/gvol0

[Install]
WantedBy=multi-user.target
